# op1/op1-server.socket
[Unit]
Description=op1-server socket

[Socket]
ListenStream=/run/op1-server.socket

[Install]
WantedBy=sockets.target
